# deploy/manibandha-go.service
# systemd unit для Go-бэкенда Manibandha.
# Установка:
#   cp deploy/manibandha-go.service /etc/systemd/system/manibandha-go.service
#   systemctl daemon-reload && systemctl enable --now manibandha-go
[Unit]
Description=Manibandha Go backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/manibandha.prema.su
# .env и uploads вынесены в корень приложения (backend/ с Python удалён).
Environment=ENV_FILE=/var/www/manibandha.prema.su/.env
Environment=PORT=8020
ExecStart=/var/www/manibandha.prema.su/mani-go
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
