# deploy/manibandha-api.service
# systemd unit for the Manibandha FastAPI backend (Gunicorn + Uvicorn workers).
# Install: cp to /etc/systemd/system/manibandha-api.service && systemctl enable --now manibandha-api
[Unit]
Description=Manibandha FastAPI backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/manibandha.prema.su/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/var/www/manibandha.prema.su/backend/.env
ExecStart=/var/www/manibandha.prema.su/backend/.venv/bin/gunicorn app.main:app \
    --workers 3 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8010 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
